# deploy/systemd/job-star-worker-gatehouse.service
[Unit]
Description=Job-Star gatehouse-ai expert worker (PR-based, machine-pinned to DESKTOP-RNK6J72)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/home/craig/job-star/job-star
Environment=PYTHONPATH=/home/craig/job-star/job-star
Environment=DATABASE_URL=postgresql://jobstar:jobstar@localhost:5432/job_star
Environment=GATEHOUSE_API_URL=http://gatehouse-ai.craigdfrench.com/v1
Environment=JOB_STAR_WORKER=gatehouse-ai-expert
Environment=JOB_STAR_EXPERT=gatehouse-ai
Environment=JOB_STAR_MACHINE=DESKTOP-RNK6J72
# Go is required for `go test ./...` in the PR executor's test loop.
Environment=PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
User=craig
Group=craig
ExecStart=/home/craig/job-star/job-star/.venv/bin/python -m job_star.worker --interval 15
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target